# systemd/sensor-logger.service
[Unit]
Description=Mock sensor logger service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sensor-logger --interval 5
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
